# supervisord-unraid-fix.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/
environment=DISPLAY=":1",HOME="/root",SHELL="/bin/bash",USER="root",LANG="en_US.UTF-8",LC_ALL="en_US.UTF-8"

[supervisorctl]
serverurl=unix:///dev/shm/supervisor.sock

[unix_http_server]
file=/dev/shm/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# PRIORITY 1: X11 Display Server
[program:xvfb]
priority=1
directory=/tmp
command=/usr/bin/Xvfb :1 -screen 0 1920x1080x24 -ac +extension GLX +render -noreset -nolisten tcp
user=root
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/xvfb.log
stderr_logfile=/var/log/xvfb.err
startsecs=2
startretries=3

# PRIORITY 2: VNC Server
[program:x11vnc]
priority=2
command=/bin/bash -c "sleep 3 && /usr/bin/x11vnc -display :1 -nopw -forever -shared -rfbport 5900 -listen 0.0.0.0 -xkb"
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/x11vnc.log
stderr_logfile=/var/log/x11vnc.err
startsecs=5
startretries=3
environment=DISPLAY=":1",HOME="/root",XAUTHORITY="/root/.Xauthority"

# PRIORITY 3: noVNC Web Interface
[program:novnc]
priority=3
command=/bin/bash -c "sleep 8 && cd /usr/share/novnc && python3 -m websockify --web . 8080 localhost:5900"
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/novnc.log
stderr_logfile=/var/log/novnc.err
startsecs=8
startretries=3
environment=DISPLAY=":1",HOME="/root"

# PRIORITY 4: XFCE Desktop Environment
[program:xfce4]
priority=4
command=/bin/bash -c "sleep 5 && /usr/bin/startxfce4"
environment=DISPLAY=":1",HOME="/root",XDG_SESSION_TYPE="x11"
user=root
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/xfce4.log
stderr_logfile=/var/log/xfce4.err
startsecs=5

# PRIORITY 5: Database Initialization
[program:database-init]
priority=5
command=/bin/bash -c "sleep 10 && /app/Entry-Database-First.sh"
autostart=true
autorestart=false
startretries=1
user=root
environment=DISPLAY=":1",
    HOME="/root",
    MYSQL_USER="farmboy",
    MYSQL_PASSWORD="Sntioi004!",
    MYSQL_ROOT_PASSWORD="Sntioi004!",
    MYSQL_HOST="mariadb",
    MYSQL_PORT="3306",
    MYSQL_DATABASE="farmboy_db"
stdout_logfile=/var/log/database-init.log
stderr_logfile=/var/log/database-init.err
startsecs=10

# PRIORITY 10: Farm Manager
[program:farm-manager]
priority=10
command=/bin/bash -c "sleep 20 && cd /app && node server.js"
user=root
autostart=true
autorestart=true
stderr_logfile=/var/log/farm-manager.err.log
stdout_logfile=/var/log/farm-manager.out.log
environment=NODE_ENV="production",
    DATABASE_URL="mysql://farmboy:Sntioi004!@mariadb:3306/farmboy_db",
    PORT="3333",
    MYSQL_USER="farmboy",
    MYSQL_PASSWORD="Sntioi004!",
    MYSQL_HOST="mariadb",
    MYSQL_PORT="3306",
    MYSQL_DATABASE="farmboy_db"
startsecs=10
startretries=3

# PRIORITY 30: EternalFarm Agent
[program:eternalfarm-agent]
priority=30
command=/bin/bash -c "sleep 45 && cd /root && DISPLAY=:1 xfce4-terminal --hold --title='EternalFarm Agent' --command='/usr/local/bin/EternalFarmAgent --auto-login --key-file=/root/EternalFarm/Logs/agent.key --show-gui'"
autostart=true
autorestart=true
startretries=5
user=root
environment=DISPLAY=":1",
    HOME="/root",
    XAUTHORITY="/root/.Xauthority",
    JAVA_HOME="/usr/lib/jvm/temurin-8-jdk-amd64"
stdout_logfile=/var/log/eternalfarm-agent.log
stderr_logfile=/var/log/eternalfarm-agent.err
startsecs=10

# PRIORITY 40: EternalFarm Checker
[program:eternalfarm-checker]
priority=40
command=/bin/bash -c "sleep 60 && cd /root && /usr/local/bin/EternalFarmChecker --auto-login --key-file=/root/EternalFarm/Logs/checker.key --headless"
autostart=true
autorestart=true
startretries=5
user=root
environment=DISPLAY=":1",
    HOME="/root",
    XAUTHORITY="/root/.Xauthority",
    JAVA_HOME="/usr/lib/jvm/temurin-8-jdk-amd64"
stdout_logfile=/var/log/eternalfarm-checker.log
stderr_logfile=/var/log/eternalfarm-checker.err
startsecs=10

# PRIORITY 50: EternalFarm Browser Automator
[program:eternalfarm-automator]
priority=50
command=/bin/bash -c "sleep 75 && cd /root && /usr/local/bin/EternalFarmBrowserAutomator --auto-login --key-file=/root/EternalFarm/Logs/automator.key --headless"
autostart=true
autorestart=true
startretries=5
user=root
environment=DISPLAY=":1",
    HOME="/root",
    XAUTHORITY="/root/.Xauthority",
    JAVA_HOME="/usr/lib/jvm/temurin-8-jdk-amd64"
stdout_logfile=/var/log/eternalfarm-automator.log
stderr_logfile=/var/log/eternalfarm-automator.err
startsecs=10 